# rc.local
#!/bin/sh -e


sleep 5

/home/pi/httpaudioplayer/start.sh &

exit 0