# rc.local
#!/bin/sh -e


sleep 5

cd /root/simpleplayer

screen -AdmS player -t home
screen -S player -X screen -t init /root/simpleplayer/start.sh

exit 0